# guard-service/akhai-guard.service
[Unit]
Description=AkhAI Guard (NLI factuality)
After=network.target

[Service]
User=guard
WorkingDirectory=/opt/akhai-guard
Environment=GUARD_TOKEN=CHANGE_ME
ExecStart=/opt/akhai-guard/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8731 --workers 1
Restart=always
RestartSec=3
MemoryMax=1500M

[Install]
WantedBy=multi-user.target
